# installer/linux/ranni.service
[Unit]
Description=RANNI - Asistente Virtual de Escritorio
After=network.target sound.target
Wants=ollama.service

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/python3 /opt/ranni/core/main.py
WorkingDirectory=/opt/ranni
Restart=always
RestartSec=5
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/var/log/ranni/ranni.log
StandardError=append:/var/log/ranni/ranni-error.log

[Install]
WantedBy=default.target
